# packaging/linux/ultrasqld.service
[Unit]
Description=UltraSQL SQL database server
Documentation=https://docs.ultrasql.org/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ultrasql
Group=ultrasql
EnvironmentFile=-/etc/ultrasql/ultrasqld.env
# All flags are taken from the environment file (ultrasqld reads ULTRASQL_LISTEN,
# ULTRASQL_DATA_DIR, ULTRASQL_AUTH_*, ULTRASQL_HBA_FILE, ULTRASQL_TLS_*, etc.).
# The default bind is loopback-only; exposing it requires configuring auth.
ExecStart=/usr/bin/ultrasqld
Restart=on-failure
RestartSec=2s
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/ultrasql
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
